# etc/crontab
#### dehydrated
0 3,15 * * * root /root/dehydrated -f /dehydrated/config --cron > /proc/1/fd/1 2> /proc/1/fd/2
