# etc/services.d/bind9/run
#!/usr/bin/with-contenv bash

exec named -4 -f
